# recipes-pironex/pia-am335x/pia-gprshelper/pppmonitor.service
[Unit]
Description=PPP Connection Test

[Service]
Type=simple
ExecStart=/usr/sbin/ppp_connection_check.sh
