# recipes-pironex/pia-am335x/pia-gprshelper/pppmonitor.service
[Unit]
Description=Test for internet connection and restart pppd if needed

[Service]
Type=simple
ExecStart=/usr/sbin/ppp_connection_check.sh
